# Platforms/Linux/BBSd
#!/bin/sh
### BEGIN INIT INFO
# Provides:          BBSd 
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Should-Start:      $syslog
# Should-Stop:       $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts or stops the Casasoft BBS daemon.
# Description:       Starts or stops the Casasoft BBS daemon.
### END INIT INFO

# clear poisonned environment
unset TMPDIR

NAME=BBSd
DAEMON=/usr/sbin/$NAME
PIDFILE=/run/$NAME.pid

test -x "$DAEMON" || exit 0

. /lib/lsb/init-functions

case "$1" in
    start)
        log_daemon_msg "Starting Casasoft BBS daemon " "$NAME"
        start-stop-daemon --make-pidfile --pidfile "$PIDFILE" --start --quiet --background --exec "$DAEMON" 
        log_end_msg $?
        ;;
    stop)
        log_daemon_msg "Stopping Casasoft BBS daemon " "$NAME"
        start-stop-daemon --pidfile "$PIDFILE" --stop --signal 3 --quiet --oknodo --exec "$DAEMON"
		rm $PIDFILE
        log_end_msg $?
        ;;
    restart|force-reload)
        $0 stop
        $0 start
        ;;
    *)
        echo "Usage: /etc/init.d/$NAME {start|stop|force-reload|restart}"
        exit 1
        ;;
esac

exit 0
